# remaster/chroot/sound/hardware/sound/scripts/bytcr_rt5651_sound.service
[Unit]
After=alsa-state.service

[Service]
Type=simple
ExecStart=/root/hardware/sound/setup_bytcr_rt5651_asound_state.sh
