# remaster/chroot/sound/hardware/sound/scripts/bytcr_rt5651_sound.service
[Unit]
After=sound.target
After=alsa-restore.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/root/hardware/sound/setup_bytcr_rt5651_asound_state.sh

[Install]
WantedBy=multi-user.target
