# remaster/chroot/sound/hardware/sound/scripts/bytcr_rt5651_sound.service
[Unit]

[Service]
ExecStart=/root/hardware/sound/setup_bytcr_rt5651_asound_state.sh
RemainAfterExit=yes

[Install]
WantedBy=alsa-state.service
