# hardware/etc_udev_rules.d/bluetooth_r8723bs.rules
SUBSYSTEM=="rfkill", RUN+="/root/hardware/scripts/bluetooth_udev_rules.sh hci0"
